# overlay-debian/etc/profile.d/dae-gateway.sh
printf '\nRun gateway to open the Debian daed Gateway menu.\n'
