# overlay-debian/etc/systemd/system/daed.service.d/10-gateway.conf
[Unit]
After=dae-gateway-firstboot.service check-ebpf.service paopaodns.service
Requires=check-ebpf.service

[Service]
Environment=HOME=/root
Environment=XDG_DATA_HOME=/root/.local/share
