# systemd/late-bind.service
[Unit]
Description=Late-bind LAN IP for daed/mosdns
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/custom-services/scripts/late-bind.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
