# overlay-debian/etc/systemd/system/ssh.service.d/10-dae-gateway.conf
[Unit]
After=dae-gateway-firstboot.service dae-ssh-hostkeys.service
Wants=dae-ssh-hostkeys.service

[Service]
ExecStartPre=/usr/bin/ssh-keygen -A
