# overlay-debian/etc/systemd/system/dae-gateway-firstboot.service
[Unit]
Description=dae gateway firstboot wizard
After=local-fs.target systemd-networkd.service
Before=getty@tty1.service
ConditionPathExists=!/etc/dae-gateway-firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dae-gateway-firstboot
StandardInput=tty-force
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
