# overlay-debian/etc/systemd/system/dae-ssh-hostkeys.service
[Unit]
Description=Generate SSH host keys on first boot
DefaultDependencies=no
Before=ssh.service
ConditionPathExistsGlob=!/etc/ssh/ssh_host_*_key

[Service]
Type=oneshot
ExecStart=/usr/bin/ssh-keygen -A

[Install]
WantedBy=multi-user.target
