# overlay-debian/etc/systemd/system/dae-qos.service
[Unit]
Description=dae gateway QoS restore
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/qos-manager apply
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
