# overlay-debian/etc/systemd/system/dns-leak-guard.service
[Unit]
Description=Apply DNS leak guard to daed routing
After=network-online.target gateway-network-init.service
Before=daed.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=DNS_LEAK_GUARD_NO_RESTART=1
ExecStart=/usr/local/sbin/dns-leak-guard apply
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
