# overlay-debian/etc/systemd/system/gateway-cache-warm.timer
[Unit]
Description=Run gateway cache warmup periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=30s
Unit=gateway-cache-warm.service

[Install]
WantedBy=timers.target
