# overlay-debian/etc/systemd/system/paopaodns.service
[Unit]
Description=PaoPaoDNS container
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
Conflicts=mini-ppdns.service

[Service]
Type=simple
LimitNOFILE=1048576
EnvironmentFile=-/etc/paopaodns/paopaodns.env
ExecStartPre=/usr/local/sbin/paopaodns-load-image
ExecStartPre=-/usr/bin/docker rm -f paopaodns
ExecStart=/usr/bin/docker run --rm --name paopaodns --network host --env-file /etc/paopaodns/paopaodns.env -v /var/lib/paopaodns:/data ${PAOPAODNS_IMAGE}
ExecStop=/usr/bin/docker stop paopaodns
ExecStopPost=-/usr/bin/docker rm -f paopaodns
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
