# overlay/etc/init.d/dae-qos
#!/sbin/openrc-run

description="daed gateway CAKE/IFB QoS service"
command="/usr/local/sbin/qos-manager"
command_args="apply"

depend() {
    need net
    before daed
}

stop() {
    ebegin "Clearing daed gateway QoS"
    /usr/local/sbin/qos-manager clear
    eend $?
}
