# overlay/etc/init.d/daed
#!/sbin/openrc-run

description="daed dashboard and transparent proxy service"
command="/usr/bin/daed"
command_args="run -c /etc/daed/"
command_background="yes"
pidfile="/run/daed.pid"
output_log="/var/log/daed/daed.log"
error_log="/var/log/daed/daed.err"
start_stop_daemon_args="--env XDG_DATA_HOME=/etc"

depend() {
    need net check-ebpf
    after mini-ppdns firewall
}

start_pre() {
    checkpath --directory --mode 0755 /etc/daed
    checkpath --directory --mode 0755 /var/log/daed
    if [ ! -x /usr/bin/daed ]; then
        eerror "/usr/bin/daed is missing or not executable"
        return 1
    fi
    if [ ! -r /etc/daed/geoip.dat ] || [ ! -r /etc/daed/geosite.dat ]; then
        eerror "geoip.dat and geosite.dat must exist in /etc/daed"
        return 1
    fi
}
