# overlay/etc/init.d/daed-firstboot
#!/sbin/openrc-run

description="daed gateway firstboot wizard"
command="/usr/local/sbin/daed-firstboot"
command_args="run"
keyword -timeout

depend() {
    need localmount
    after bootmisc gateway-network-init
    before daed sshd
}

start() {
    if [ -f /etc/dae-gateway-firstboot.done ]; then
        ebegin "Skipping daed gateway firstboot"
        eend 0
        return 0
    fi

    ebegin "Running daed gateway firstboot"
    /usr/local/sbin/daed-firstboot run </dev/console >/dev/console 2>&1
    eend $?
}
