# overlay/etc/profile.d/dae-gateway.sh
#!/bin/sh

if [ -n "${PS1:-}" ]; then
    printf '\nRun gateway to open the daed Alpine Gateway menu.\n'
fi
